# systemd/dc-reload@.service
[Unit]
Description=Refresh %i images and update containers

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload-or-restart dc@%i
